# templates/jira-init.erb
#! /bin/sh

### BEGIN INIT INFO
# Provides:          jira
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO

. /lib/lsb/init-functions

NAME="jira"
JIRA_HOME="<%= scope.lookupvar('jira::params::home') %>"
JIRA_PATH="/opt/atlassian-jira-<%= scope.lookupvar('jira::params::version') %>-standalone"
START_SCRIPT="${JIRA_PATH}/bin/startup.sh"
STOP_SCRIPT="${JIRA_PATH}/bin/stop_jira.sh"
PIDFILE="${JIRA_PATH}/work/catalina.pid"
CWD=`pwd`

start () {
        log_daemon_msg "Starting $NAME"
        if start-stop-daemon --start --quiet --user jira -d $JIRA_PATH --pidfile "$PIDFILE" --exec $START_SCRIPT; then
                log_end_msg 0
        else
                log_end_msg 1
        fi
}

stop () {
        start-stop-daemon --stop --quiet --oknodo --pidfile "$PIDFILE" --exec $STOP_SCRIPT
}

status () {
        status_of_proc -p $PIDFILE "" "$NAME"
}

case $1 in
        start)
                if status; then exit 0; fi
                start
                ;;
        stop)
                stop
                ;;
        reload)
                stop
                sleep 2
                start
                ;;
        restart)
                stop
                sleep 2
                start
                ;;
        status)
                status && exit 0 || exit $?
                ;;
        *)
                echo "Usage: $0 {start|stop|restart|reload|status}"
                exit 1
                ;;
esac

exit 0
